# crontab/crontab-p
# Edit this file to introduce tasks to be run by cron.
#
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
#
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
#
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
#
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
#
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
#
# For more information see the manual pages of crontab(5) and cron(8)
#

MAILTO="wer987wer888@gmail.com"

# will run in all non-commmented out lines manually
# crontab -l | grep -v '^\s*#' | sed -E 's/.*cronic\s+//'  # lists
# crontab -l | grep -v '^\s*#' | sed -E 's/.*cronic\s+//' | while read -r cmd; do bash -c "$cmd" || true; done

# @daily, @weekly, @monthly, @hourly, @yearly, @midnight (alias for @daily)
#@hourly    → minute 0 of every hour            (e.g. 01:00, 02:00, ...)
#@daily     → 00:00 every day                   (midnight)
#@weekly    → 00:00 every Sunday                (midnight at the start of Sunday)
#@monthly   → 00:00 on the 1st of every month
#@yearly    → 00:00 on January 1st              (same as @annually)

# m h dom mon dow   command
0 2 * * * cronic /home/user/dotfiles/jo_alerts 
@weekly cronic /home/user/dotfiles/email_test         # every Sunday at 00:00 (weekly)
@hourly cronic /home/user/dotfiles/check-space        # every hour on the hour
0 2 * * * cronic /home/user/dotfiles/snap_clean         # daily at 02:00
0 4 * * * cronic /home/user/dotfiles/clamscan_daily         # daily at 04:00
@daily cronic cd /home/user/dotfiles && git pull # --rebase && git push # bad idea? dont use --rebase bc "replay" local over origin

